# data/ayatana-indicator-display.service.in
[Unit]
Description=Ayatana Indicator Display Backend
PartOf=graphical-session.target
After=ayatana-indicators-pre.target

[Service]
ExecStart=@pkglibexecdir@/ayatana-indicator-display-service
Restart=on-failure
