# data/ayatana-indicator-display.service.in
[Unit]
Description=Indicator Display Backend
PartOf=graphical-session.target
After=indicators-pre.target

[Service]
ExecStart=@pkglibexecdir@/indicator-display-service
Restart=on-failure
