# data/ayatana-indicator-display.conf.in
description "Ayatana Indicator Display Backend"

start on ayatana-indicator-services-start
stop on desktop-end or ayatana-indicator-services-end

respawn
respawn limit 2 10

exec @CMAKE_INSTALL_FULL_LIBEXECDIR@/ayatana-indicator-display/ayatana-indicator-display-service
