# data/ayatana-indicator-display.service.in
[Unit]
Description=Ayatana Indicator Display Backend
PartOf=graphical-session.target
PartOf=ayatana-indicators.target

[Service]
ExecStart=@pkglibexecdir@/ayatana-indicator-display-service
Restart=on-failure

[Install]
WantedBy=ayatana-indicators.target
